# systemd/ntn-frontend.service.template
[Unit]
Description=Not The News - Local Frontend (Vite)
After=network.target ntn-backend.service

[Service]
Type=simple
WorkingDirectory={{REPO_DIR}}
ExecStart=/usr/bin/npx vite --mode development
Restart=always
Environment=PATH=/usr/bin:/usr/local/bin

[Install]
WantedBy=default.target
