# systemd/ntn-backend.service
[Unit]
Description=Not The News - Local Backend (Wrangler)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ntn/git/not-the-news/worker
ExecStart=/usr/bin/npx wrangler dev
Restart=always
Environment=PATH=/usr/bin:/usr/local/bin

[Install]
WantedBy=default.target
